# spark-jobserver/roles/spark-job-server/templates/supervisord.conf.j2
[supervisord]
nodaemon=true
loglevel=debug

[program:jobserver]
redirect_stderr=true
command=/bin/bash -c /spark-jobserver/runner/server_start.sh